# configs/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/sbin/service nginx start
redirect_stderr=true

[include]
files = /etc/supervisord.d/*.ini
